# java/oraclejdk23/profile.d/jdk.csh
#!/bin/csh

setenv JAVA_HOME /usr/lib${LIBDIRSUFFIX}/java
setenv PATH \${PATH}:\${JAVA_HOME}/bin
